# ExpertAgentC_LLMServer_Nursing_Note_STT/gemma-audio-api/gemma-api.service
[Unit]
Description=Gemma-3n-E4B Audio API Service
After=network.target

[Service]
Type=simple
User=benchmark
WorkingDirectory=/home/benchmark/google_hackathon/gemma-audio-api
Environment=PATH=/home/benchmark/google_hackathon/gemma-audio-api/venv/bin
ExecStart=/home/benchmark/google_hackathon/gemma-audio-api/venv/bin/python /home/benchmark/google_hackathon/gemma-audio-api/gemma_api.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
